# configs/etc/cron.d/photonbbs
*/5 * * * * chat /opt/photonbbs/bin/maint.sh >/dev/null 2>&1
00 00 * * * chat mv -f /opt/photonbbs/data/text/lastcallers /opt/photonbbs/data/text/lastcallers.last >/dev/null 2>&1
