# configs/etc/xinetd.d/photonbbs
service photonbbs
{
        disable = no
        flags           = REUSE
        socket_type     = stream
        per_source      = 1
        wait            = no
        user            = bbs
        server          = /usr/sbin/in.telnetd
        server_args     = -N -h -n -L /opt/photonbbs/bbs.pl
        log_on_failure  += USERID
}
